# fetch-image-worker.service.j2
[Unit]
Description=Fetch Image Processing Worker
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/fetch-app
# Use uv run to ensure the script uses the project's virtual environment
ExecStart=/root/.local/bin/uv run python3 /opt/fetch-app/scripts/process_images.py
